# config/service/onedrive.conf
[program:onedrive]
command=/usr/bin/onedrive.sh
process_name=onedrive
#numprocs=1
autostart=false
autorestart=true
startretries=999
redirect_stderr=true